# alpine/crontab/apache
0       1       *       *       *       /usr/bin/php84 /var/www/localhost/htdocs/cron.php [key]
*       *       *       *       *       /usr/bin/php84 /var/www/localhost/htdocs/cron_ticket_email_parser.php [key]
*       *       *       *       *       /usr/bin/php84 /var/www/localhost/htdocs/cron_mail_queue.php [key]
0       2       *       *       *       /usr/bin/php84 /var/www/localhost/htdocs/cron_certificate_refresher.php [key]
0       3       *       *       *       /usr/bin/php84 /var/www/localhost/htdocs/cron_domain_refresher.php [key]
